# packages/base/any/mgmt/builds/units/gs-mgmt.target
[Unit]
Description=Goldstone Management Daemons

[Install]
WantedBy=multi-user.target
